# playbook/daemon/tg_bot_factorio.service.j2
[Unit]
Description=TG bot Factorio

[Service]
Restart=on-failure
RestartSec=5s
WorkingDirectory=${PWD}
ExecStart=python3 ${PWD}/py_files/main.py